# airtime-silan.service
[Unit]
Description=Airtime Silan silence detector

[Service]
Type=simple
ExecStart=/usr/bin/airtime-silan
